# com.binreminder.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.binreminder.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/robertfluellen/Documents/claudework/bin-reminder/bin_reminder.py</string>
    </array>

    <!-- Run every day at 7:00 PM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/robertfluellen/Documents/claudework/bin-reminder/logs/output.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/robertfluellen/Documents/claudework/bin-reminder/logs/error.log</string>
</dict>
</plist>
